# systemd/user/fig-daemon.service
[Unit]
Description=Fig Daemon

[Service]
ExecStart=fig daemon
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
